# project/main.cron
* 2 * * * /home/de11an/kore/project/main.py &> /home/de11an/kore/project/log.txt